# 1/mydependent.service
[Unit]
Description=My Dependent Service
After=network.target
Wants=network.target

[Service]
ExecStart=/bin/bash -c 'sleep infinity'

[Install]
WantedBy=default.target
