# 1/mytimer.service
[Unit]
Description=My Timer Service

[Service]
ExecStart=/usr/bin/bash -c '/usr/bin/echo "My Timer Service ran at $(date)" >> ~/mytimer.log'
